# systemd/admin-bot.service
[Unit]
Description=TG Auto-Reply Admin Bot
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=aiuser
Group=aiuser
WorkingDirectory=/home/aiuser/tg-auto-reply
Environment=PATH=/home/aiuser/tg-auto-reply/venv/bin
EnvironmentFile=/home/aiuser/tg-auto-reply/.env
ExecStart=/home/aiuser/tg-auto-reply/venv/bin/python -m bots.admin_bot
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
